# packaging/archlinux/gandictl/gandictl-ddns.timer
[Unit]
Description=Hourly Gandi Dynamic DNS Updater

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target